# etc/supd/lserver.applet.cms.conf
[program:lserver.applet.cms]
depends_on=lserver.0.beanstald,lserver.0.redis
user=root
environment=PRJ_ROOT="/app",GIN_MODE="release",LD_LIBRARY_PATH=""
command=/app/src/applet/cms/cms
autostart=true
autorestart=true
directory=/app/src/applet/cms
stdout_logfile=/app/var/log/lserver.applet.cms.logfile.stdout
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stdout_capture=/app/var/log/lserver.applet.cms.captuer.stdout
stdout_capture_maxbytes=10MB
stdout_capture_backups=10
stderr_logfile=/app/var/log/lserver.applet.cms.logfile.stderr
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=10
stderr_capture=/app/var/log/lserver.applet.cms.capture.stderr
stderr_capture_maxbytes=10MB
stderr_capture_backups=10

